# scripts/com.milkroad.macro-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.milkroad.macro-refresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__REPO__/scripts/refresh-if-et-time.sh</string>
        <string>briefs</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__HOME__/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- Weekly brief refresh: Monday only, 5 minutes after the data refresh,
         so the brief reflects the freshly-pulled Monday close data.
         launchd uses Prague local time, so fire at both possible Prague
         equivalents of 4:05pm ET and let refresh-if-et-time.sh run only
         when America/New_York is Mon 16:05. Weekday 1=Mon. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>5</integer></dict>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>22</integer><key>Minute</key><integer>5</integer></dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__REPO__/.cache/launchd-refresh.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO__/.cache/launchd-refresh.log</string>
</dict>
</plist>
